# etc/containers/systemd/hello-world.container
[Container]
Image=docker.io/library/hello-world:latest
ContainerName=hello-world